# etc/sv/ussd-proxy/run
#!/bin/bash

. /etc/ussd-proxy.config
exec 2>&1
echo "ussd-proxy restart" | /usr/bin/ts >> /var/log/runsv.log
exec ussd-proxy $COMMAND_LINE
