# etc/sv/reg-proxy/run
#!/bin/bash

. /etc/reg-proxy.config
exec 2>&1
echo "reg-proxy restart" | /usr/bin/ts >> /var/log/runsv.log
exec reg-proxy $COMMAND_LINE
